# scripts/openflight.service
[Unit]
Description=OpenFlight Golf Launch Monitor
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=coleman
WorkingDirectory=/home/coleman/openflight
Environment=DISPLAY=:0
ExecStart=/home/coleman/openflight/scripts/start-kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
